# .config/systemd/user/docker-compose-reload.service
[Unit]
Description=Refresh images and update containers

[Service]
Type=oneshot

ExecStart=/bin/systemctl --user reload-or-restart docker-compose.service
