# .config/systemd/user/generate-lockscreen.service
[Unit]
Description=Generate a new lockscreen image

[Service]
Type=oneshot

ExecStart=/home/luzifer/bin/generate_lockscreen
