# .config/systemd/user/lock.service
[Unit]
Description=Lock screen and stop services

[Service]
ExecStart=%h/bin/lock
