# .config/systemd/user/ssh-agent.service
[Service]
ExecStart=/usr/bin/ssh-agent -D -a "${HOME}/.ssh/ssh_auth_sock"

[Install]
WantedBy=default.target
